# claude_code_messages/rootfs/etc/services.d/backend/run
#!/usr/bin/with-contenv bashio
# Launch the FastAPI backend
cd /app
export PORT="$(bashio::config 'PORT' '8099')"
exec python3 -u server.py
