# claude_code_messages/rootfs/etc/services.d/backend/run
#!/usr/bin/with-contenv bashio
# Launch the FastAPI backend
cd /app
exec python3 -u server.py
